# spiffer/99-spiffer.rules
## spiffer events
##
## Inivation cameras
ACTION=="add", ATTR{idVendor}=="152a", RUN+="/usr/bin/pkill --signal SIGUSR1 spiffer"
